# cinder/files/cinder-device-setup.conf
description "Create a loopback device for LVM Cinder volume"

start on filesystem
task

exec losetup /dev/loop0 {{ local_file }}
